# deploy/worldcup.service
# systemd unit for the World Cup prediction game.
# Install to /etc/systemd/system/worldcup.service
#
# Assumes:
#   - app cloned at /home/ubuntu/easygds_worldcup
#   - Node installed system-wide (NodeSource) so node_modules/.bin/next resolves
#   - env file at /etc/worldcup.env
#
# After editing paths/user, run:
#   sudo systemctl daemon-reload && sudo systemctl enable --now worldcup

[Unit]
Description=World Cup 2026 prediction game (Next.js)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/easygds_worldcup
EnvironmentFile=/etc/worldcup.env
ExecStart=/home/ubuntu/easygds_worldcup/node_modules/.bin/next start
Restart=on-failure
RestartSec=5

# Let the app write only where it needs to.
ReadWritePaths=/var/lib/worldcup
StateDirectory=worldcup
NoNewPrivileges=true
ProtectSystem=full

[Install]
WantedBy=multi-user.target
